# raspberry/allskySCRIPT/crontab
@reboot /home/pi/SADR/raspberry/allskySCRIPT/startMETEOLLSKY.sh  >> /media/data/crontab.log 2>&1
